# recipes-graphics/toyota/files/homescreen.service.in
[Unit]
Requires=@Requires@
After=@After@

[Service]
User=@User@
Group=@Group@
Environment=@Environment@
ExecStart=@ExecStart@

[Install]
WantedBy=@WantedBy@
