# recipes-graphics/toyota/files/homescreen.service.in
[Unit]
@ServiceUnit@

[Service]
@ServiceService@

[Install]
@ServiceInstall@
